# git-push.service
[Unit]
Description=Git Push Service

[Service]
ExecStart=/opt/tiger/PFLlib/git-push.sh

[Install]
WantedBy=default.target
